# etc/sv/installer/run
#!/bin/bash

/opt/void-autoinstaller/install.sh
