# etc/sv/installer/run
#!/bin/bash

/opt/void-autoinstall/install.sh
